# fqdn-scanner.service
[Unit]
Description=FQDN Security Scanner Web UI
After=network.target

[Service]
Type=simple
User=azureuser
WorkingDirectory=/home/azureuser/.openclaw/workspace/fqdn-security-scanner
Environment="PATH=/home/azureuser/.openclaw/workspace/fqdn-security-scanner/venv/bin"
ExecStart=/home/azureuser/.openclaw/workspace/fqdn-security-scanner/venv/bin/python3 -m flask --app web.app run --host 127.0.0.1 --port 5000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
